# github-runner/rootfs/etc/services.d/runner/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the GitHub Actions runner.
#
# Reads options from /data/options.json (populated by Supervisor from the
# add-on's Options tab).
#
# Credential persistence model (v1.0.3+):
#   - The actions/runner binary writes its long-lived credentials to
#     /opt/runner/.credentials, .credentials_rsaparams, and .runner after
#     successful first-time registration.
#   - HA Supervisor recreates the add-on container on Stop/Start (NOT just
#     stops/starts the same one). The /opt/runner/ writable layer is wiped.
#     Only the /data mount survives container recreation.
#   - So we MIRROR the three credential files to /data/runner-credentials/
#     after first registration, and RESTORE them from there on every start
#     before the existence check. This makes Stop/Start cycles work without
#     a fresh registration token.
#
# When IS a fresh token still needed?
#   - First-ever start after install (no credentials yet, anywhere).
#   - If you manually delete /data/runner-credentials/.
#   - If the runner was manually removed in the GitHub UI (the credentials
#     are revoked server-side; they'll fail auth on next connect).
# ==============================================================================

# The actions/runner binary refuses to start as root by default
# ("Must not run with sudo"). This is the right default for bare-metal
# installs, but every HA add-on container runs as root by design —
# isolation comes from the container, not from a user account.
# RUNNER_ALLOW_RUNASROOT=1 is the official escape hatch from the
# actions/runner project for this exact situation.
export RUNNER_ALLOW_RUNASROOT="1"

REPO=$(bashio::config 'github_repo')
TOKEN=$(bashio::config 'github_runner_token')
RUNNER_NAME=$(bashio::config 'runner_name')
LABELS=$(bashio::config 'runner_labels')
WORKDIR=$(bashio::config 'runner_workdir')

PERSIST_DIR="/data/runner-credentials"
mkdir -p "${WORKDIR}" "${PERSIST_DIR}"
cd /opt/runner || exit 1

# Restore credential files from the persistent /data mirror, if present.
# Three files together form the runner's identity:
#   .credentials           — bearer token used to connect to GitHub
#   .credentials_rsaparams — RSA key material for message signing
#   .runner                — runner-side config (id, name, labels, work dir)
# All three must move as a set, or auth fails with corrupt-credentials errors.
if [[ -f "${PERSIST_DIR}/.credentials" && ! -f /opt/runner/.credentials ]]; then
  cp -p "${PERSIST_DIR}/.credentials"           /opt/runner/.credentials
  cp -p "${PERSIST_DIR}/.credentials_rsaparams" /opt/runner/.credentials_rsaparams
  cp -p "${PERSIST_DIR}/.runner"                /opt/runner/.runner
  bashio::log.info "Restored runner credentials from /data — no token required."
fi

if [[ -f /opt/runner/.credentials ]]; then
  bashio::log.info "Runner credentials present — skipping registration."
else
  if bashio::var.is_empty "${TOKEN}"; then
    bashio::log.fatal "First-time registration needs a runner token."
    bashio::log.fatal "Generate at: Repo → Settings → Actions → Runners →"
    bashio::log.fatal "New self-hosted runner → Linux. Copy the value after"
    bashio::log.fatal "--token, paste into the Options tab, save, and restart."
    exit 1
  fi

  bashio::log.info "Registering runner '${RUNNER_NAME}' against ${REPO}..."
  ./config.sh \
    --url "https://github.com/${REPO}" \
    --token "${TOKEN}" \
    --name "${RUNNER_NAME}" \
    --labels "${LABELS}" \
    --work "${WORKDIR}" \
    --unattended \
    --replace \
    || bashio::exit.nok "Runner registration failed — see the lines above for the actual error from config.sh. Common causes: expired token (tokens last ~1h), wrong github_repo, network blocked."

  # Mirror the freshly-written credentials to /data so the next Stop/Start
  # (which recreates the container) doesn't need a new token.
  cp -p /opt/runner/.credentials           "${PERSIST_DIR}/.credentials"
  cp -p /opt/runner/.credentials_rsaparams "${PERSIST_DIR}/.credentials_rsaparams"
  cp -p /opt/runner/.runner                "${PERSIST_DIR}/.runner"
  bashio::log.info "Mirrored credentials to /data for restart-safety."
fi

bashio::log.info "Starting GitHub Actions runner: ${RUNNER_NAME}"
exec ./run.sh
