# github-runner/rootfs/etc/services.d/runner/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start one or more GitHub Actions runners (v2.0.0+).
#
# Config model:
#   - The `runners:` array in /data/options.json is the canonical source.
#     Each entry registers an independent runner against its named repo.
#   - Legacy single-runner fields (github_repo, runner_name, ...) are honoured
#     as a fallback when `runners:` is empty — preserves v1.x → v2.x upgrades
#     without user intervention.
#
# Per-runner layout:
#   /opt/runner/                       template installation (built into image)
#   /opt/runners/<name>/               per-instance installation (cp -r from template)
#   /data/runner-credentials/<name>/   persisted credentials (survives container recreation)
#   /data/_work/<name>/                working directory for jobs
#
# Process model:
#   Each runner runs `./run.sh` in a background subshell. The supervisor `wait`s
#   on all of them; if any exit, the supervisor exits and s6 restarts the
#   whole service group. Simple and correct for the single-NUC scale this
#   add-on targets.
#
# When IS a fresh token needed?
#   - First-ever start of a new runner entry (no /data/runner-credentials/<name>/).
#   - If you manually delete /data/runner-credentials/<name>/.
#   - If the runner was manually removed in the GitHub UI (credentials are
#     revoked server-side; they'll fail auth on next connect).
# ==============================================================================

# actions/runner refuses to start as uid=0 by default. HA add-on containers
# run as root by design (container is the isolation boundary). This env var
# is the official escape hatch from the actions/runner project.
export RUNNER_ALLOW_RUNASROOT="1"

PERSIST_BASE="/data/runner-credentials"
RUNNERS_DIR="/opt/runners"
RUNNER_TEMPLATE="/opt/runner"

# v3.0.6: hybrid cache layout.
#
# Workfolder default moves from /data/_work to /work. /work is in the
# container's writable layer — ephemeral (wiped on Stop/Start/Update) but
# NOT in /data, so HAOS backups stay small regardless of CI traffic. The
# trade-off (cold cache after container recreation) is mitigated by
# redirecting SDK + dependency caches into /data/cache/ via env vars
# below — those persist across restarts and are tiny relative to a full
# _work directory.
#
# Backwards-compat: if the user's HA option is the LEGACY default literal
# "/data/_work", treat it as "I never set this, use the new default".
# Any non-default value is honoured as an explicit user choice.
WORKDIR_BASE=$(bashio::config 'workdir_base')
if [[ -z "${WORKDIR_BASE}" || "${WORKDIR_BASE}" == "null" || "${WORKDIR_BASE}" == "/data/_work" ]]; then
  WORKDIR_BASE="/work"
fi

mkdir -p "${PERSIST_BASE}" "${WORKDIR_BASE}" "${RUNNERS_DIR}"

# Shared tool cache (RUNNER_TOOL_CACHE). subosito/flutter-action,
# actions/setup-node, actions/setup-java, etc. install SDKs here keyed
# by version, so this is safe to share across runners — different
# tools/versions get their own subdirs and don't collide.
export RUNNER_TOOL_CACHE="/data/cache/tools"
mkdir -p "${RUNNER_TOOL_CACHE}"

# ─── Build the canonical runners JSON array ──────────────────────────────────

# bashio::config returns list values as NDJSON (one element per line), not
# as a single JSON array. `jq -s` slurps the lines into an array; `map(select(. != null))`
# drops any null entries (which can appear if bashio emits a stray "null"
# line for a missing/empty key). The result is always a valid JSON array,
# possibly empty.
RUNNERS_JSON=$(bashio::config 'runners' | jq -cs 'map(select(. != null))')

RUNNER_COUNT=$(echo "${RUNNERS_JSON}" | jq 'length')

if [[ "${RUNNER_COUNT}" == "0" ]]; then
  bashio::log.fatal "No runners configured."
  bashio::log.fatal "Add at least one entry to the 'runners' list in the Options tab."
  bashio::log.fatal "Each entry needs: name, repo, labels, registration_token."
  exit 1
fi

# ─── Migrate v1.x credentials layout to v2.x ─────────────────────────────────
# v1.x stored credentials at /data/runner-credentials/.credentials (flat).
# v2.x uses /data/runner-credentials/<name>/.credentials (per-runner).
# On first v2.x start after upgrade, move the flat-layout credentials into
# the first configured runner's slot.
if [[ -f "${PERSIST_BASE}/.credentials" ]]; then
  FIRST_NAME=$(echo "${RUNNERS_JSON}" | jq -r '.[0].name')
  TARGET_DIR="${PERSIST_BASE}/${FIRST_NAME}"
  if [[ -n "${FIRST_NAME}" && ! -f "${TARGET_DIR}/.credentials" ]]; then
    mkdir -p "${TARGET_DIR}"
    mv "${PERSIST_BASE}/.credentials"           "${TARGET_DIR}/.credentials"
    mv "${PERSIST_BASE}/.credentials_rsaparams" "${TARGET_DIR}/.credentials_rsaparams" 2>/dev/null || true
    mv "${PERSIST_BASE}/.runner"                "${TARGET_DIR}/.runner" 2>/dev/null || true
    bashio::log.info "Migrated v1.x credentials to ${TARGET_DIR}/."
  fi
fi

bashio::log.info "Starting ${RUNNER_COUNT} runner(s)."

# ─── Per-runner startup ──────────────────────────────────────────────────────

start_runner() {
  local NAME="$1"
  local REPO="$2"
  local LABELS="$3"
  local TOKEN="$4"

  local INSTALL_DIR="${RUNNERS_DIR}/${NAME}"
  local PERSIST_DIR="${PERSIST_BASE}/${NAME}"
  local WORK_DIR="${WORKDIR_BASE}/${NAME}"

  mkdir -p "${PERSIST_DIR}" "${WORK_DIR}"

  # First-time setup: clone the template installation into a per-runner copy.
  # The actions/runner binary writes per-instance state (`.credentials`,
  # `.credentials_rsaparams`, `.runner`) into its install directory, so each
  # runner needs its own. Disk cost: ~600MB per runner (acceptable on HAOS
  # NUC scale).
  if [[ ! -d "${INSTALL_DIR}" ]]; then
    bashio::log.info "[${NAME}] Creating per-runner installation at ${INSTALL_DIR}..."
    cp -r "${RUNNER_TEMPLATE}" "${INSTALL_DIR}"
  fi

  # Restore credentials from /data if the container layer is fresh
  # (Stop/Start/Update recreated the container; only /data survived).
  if [[ -f "${PERSIST_DIR}/.credentials" && ! -f "${INSTALL_DIR}/.credentials" ]]; then
    cp -p "${PERSIST_DIR}/.credentials"           "${INSTALL_DIR}/.credentials"
    cp -p "${PERSIST_DIR}/.credentials_rsaparams" "${INSTALL_DIR}/.credentials_rsaparams"
    cp -p "${PERSIST_DIR}/.runner"                "${INSTALL_DIR}/.runner"
    bashio::log.info "[${NAME}] Restored credentials from /data — no token required."
  fi

  cd "${INSTALL_DIR}" || {
    bashio::log.error "[${NAME}] Could not cd to install dir."
    return 1
  }

  if [[ ! -f .credentials ]]; then
    if [[ -z "${TOKEN}" ]]; then
      bashio::log.fatal "[${NAME}] No persisted credentials and no registration token."
      bashio::log.fatal "[${NAME}] Generate a token at: ${REPO} → Settings → Actions →"
      bashio::log.fatal "[${NAME}] Runners → New self-hosted runner → Linux. Paste it into"
      bashio::log.fatal "[${NAME}] this runner entry's registration_token field."
      return 1
    fi

    bashio::log.info "[${NAME}] Registering against ${REPO}..."
    ./config.sh \
      --url "https://github.com/${REPO}" \
      --token "${TOKEN}" \
      --name "${NAME}" \
      --labels "${LABELS}" \
      --work "${WORK_DIR}" \
      --unattended \
      --replace \
      || {
        bashio::log.error "[${NAME}] Registration failed. Common causes: expired token (tokens last ~1h), wrong repo, network blocked."
        return 1
      }

    # Mirror to /data for restart-safety.
    cp -p .credentials           "${PERSIST_DIR}/.credentials"
    cp -p .credentials_rsaparams "${PERSIST_DIR}/.credentials_rsaparams"
    cp -p .runner                "${PERSIST_DIR}/.runner"
    bashio::log.info "[${NAME}] Mirrored credentials to /data."
  else
    bashio::log.info "[${NAME}] Credentials present — skipping registration."

    # v3.0.6 migration: if the persisted .runner file's workFolder points
    # at the legacy /data/_work/<name>/ path, rewrite it to the new
    # WORK_DIR. Without this, runners registered pre-v3.0.6 would keep
    # writing to /data/_work — exactly what the hybrid layout was meant
    # to avoid. jq update + atomic mv, then mirror to /data so the change
    # survives the next container recreation.
    if [[ -f .runner ]]; then
      local stored_work
      stored_work=$(jq -r '.workFolder // ""' .runner 2>/dev/null)
      if [[ -n "${stored_work}" && "${stored_work}" != "${WORK_DIR}" ]]; then
        bashio::log.info "[${NAME}] Migrating workFolder ${stored_work} → ${WORK_DIR} (v3.0.6 hybrid layout)."
        jq --arg new "${WORK_DIR}" '.workFolder = $new' .runner > .runner.tmp \
          && mv .runner.tmp .runner \
          && cp -p .runner "${PERSIST_DIR}/.runner"
      fi
    fi
  fi

  # ─── Per-runner cache env vars (v3.0.6 hybrid) ──────────────────────────
  # SDK and dependency caches stay in /data (persist across container
  # restarts, included in HAOS backups but small) while the work tree
  # itself is in /work (ephemeral, NOT in backups).
  #
  # GRADLE_USER_HOME MUST be per-runner — sharing it can deadlock when two
  # runners build Android projects simultaneously (Gradle daemon locks).
  # PUB_CACHE and npm_config_cache are version-keyed internally and could
  # be shared safely, but per-runner keeps things clean and avoids any
  # surprises if two runners end up using mismatched dep versions.
  export PUB_CACHE="/data/cache/${NAME}/pub"
  export GRADLE_USER_HOME="/data/cache/${NAME}/gradle"
  export npm_config_cache="/data/cache/${NAME}/npm"
  mkdir -p "${PUB_CACHE}" "${GRADLE_USER_HOME}" "${npm_config_cache}"

  # ─── Self-heal loop ─────────────────────────────────────────────────────
  # Wraps ./run.sh so we can detect a specific class of failure (server-side
  # registration invalidated → "registration has been deleted" or
  # "PublicAccess … not authorized") and recover automatically without manual
  # shell access to /data.
  #
  # When detected:
  #   1. Nuke this runner's persisted credentials (/data/runner-credentials/<name>/).
  #   2. Re-register against ${REPO} using whatever token is in HA Options.
  #      - If the token is still valid (rare — single-use tokens are usually
  #        already consumed by an earlier registration), this fully self-heals.
  #      - If the token is consumed/expired, re-registration fails fast and
  #        the subshell exits with a clear actionable error. The stale creds
  #        are still gone, so the next restart-with-fresh-token recovers
  #        cleanly without any manual `rm -rf` step.
  #   3. Cap retries at MAX_HEAL_ATTEMPTS to avoid infinite loops.
  #
  # Non-auth exits (clean shutdown, crash, etc.) pass through as before —
  # supervisor `wait` notices and s6 restarts the service group.
  local MAX_HEAL_ATTEMPTS=3
  local attempt=0
  while :; do
    bashio::log.info "[${NAME}] Starting runner..."
    local LOG_TMP
    LOG_TMP=$(mktemp)
    # tee passes output to the container log AND captures it for pattern matching.
    # PIPESTATUS[0] is ./run.sh's exit code (tee always returns 0 unless its
    # tmpfile fills).
    ./run.sh 2>&1 | tee "${LOG_TMP}"

    if ! grep -qE "registration has been deleted from the server|System:PublicAccess.*not authorized" "${LOG_TMP}"; then
      rm -f "${LOG_TMP}"
      break  # clean exit (or non-auth failure) — let the supervisor handle it
    fi
    rm -f "${LOG_TMP}"

    attempt=$((attempt + 1))
    if [[ ${attempt} -gt ${MAX_HEAL_ATTEMPTS} ]]; then
      bashio::log.error "[${NAME}] Self-heal exhausted ${MAX_HEAL_ATTEMPTS} attempts. Giving up. Paste a fresh registration_token in HA Options and restart the add-on."
      return 1
    fi

    bashio::log.warning "[${NAME}] Stale server-side registration detected (attempt ${attempt}/${MAX_HEAL_ATTEMPTS}). Nuking local credentials and re-registering."
    rm -f .credentials .credentials_rsaparams .runner
    rm -rf "${PERSIST_DIR}"
    mkdir -p "${PERSIST_DIR}"

    if [[ -z "${TOKEN}" ]]; then
      bashio::log.error "[${NAME}] Cannot self-heal: no registration_token set in HA Options. Stale credentials have been cleared; paste a fresh token and restart the add-on."
      return 1
    fi

    bashio::log.info "[${NAME}] Re-registering against ${REPO}..."
    if ! ./config.sh \
      --url "https://github.com/${REPO}" \
      --token "${TOKEN}" \
      --name "${NAME}" \
      --labels "${LABELS}" \
      --work "${WORK_DIR}" \
      --unattended \
      --replace; then
      bashio::log.error "[${NAME}] Re-registration failed. The token in HA Options is likely already consumed (registration tokens are single-use). Stale credentials have been cleared; paste a fresh token and restart the add-on."
      return 1
    fi

    cp -p .credentials           "${PERSIST_DIR}/.credentials"
    cp -p .credentials_rsaparams "${PERSIST_DIR}/.credentials_rsaparams"
    cp -p .runner                "${PERSIST_DIR}/.runner"
    bashio::log.info "[${NAME}] Self-heal succeeded. Mirrored fresh credentials to /data."
  done
}

# ─── Spawn all runners and wait ──────────────────────────────────────────────

# Use a tempfile rather than process substitution + while-read so the for-loop
# can use `&` (which doesn't work cleanly inside `while read` pipelines).
TMP_RUNNERS=$(mktemp)
echo "${RUNNERS_JSON}" | jq -c '.[]' > "${TMP_RUNNERS}"

# Stagger startup (v3.0.4). Without a delay, multiple runners reconnect to
# GitHub within milliseconds of each other on every restart. Symptom: one
# random runner per restart has its server-side registration silently
# invalidated and fails with `PublicAccess … not authorized` on the next
# connect. A small gap lets each runner fully establish its session before
# the next one starts. 5s is comfortably above the ~1s the runner takes to
# reach `Listening for Jobs` and well below the human-visible startup
# threshold (15s total for 3 runners).
STAGGER_SECONDS=5

# Track child PIDs so the SIGTERM/SIGINT handler can forward signals.
CHILD_PIDS=()

first=1
while IFS= read -r runner; do
  NAME=$(echo "${runner}" | jq -r '.name')
  REPO=$(echo "${runner}" | jq -r '.repo')
  LABELS=$(echo "${runner}" | jq -r '.labels')
  TOKEN=$(echo "${runner}" | jq -r '.registration_token // ""')

  if [[ -z "${NAME}" || -z "${REPO}" || -z "${LABELS}" || "${NAME}" == "null" ]]; then
    bashio::log.error "Skipping malformed runner entry: ${runner}"
    continue
  fi

  # v3.0.7 defensive validation. The name + repo values flow into filesystem
  # paths (PERSIST_DIR, INSTALL_DIR, WORK_DIR, cache dirs) and into the
  # config.sh CLI. Bash double-quoting handles arg-splitting, but a name
  # containing `../` would still resolve to a path outside the intended
  # directory tree (and `rm -rf "${PERSIST_DIR}"` in the self-heal loop
  # would then nuke an unintended location). Reject filesystem-unsafe
  # characters at config-load time, before any path is constructed.
  if [[ ! "${NAME}" =~ ^[A-Za-z0-9_-]+$ ]]; then
    bashio::log.error "Skipping runner: name '${NAME}' must match ^[A-Za-z0-9_-]+$ (filesystem-safe alphanumerics, underscore, hyphen)."
    continue
  fi
  if [[ ! "${REPO}" =~ ^[A-Za-z0-9_.-]+/[A-Za-z0-9_.-]+$ ]]; then
    bashio::log.error "Skipping runner: repo '${REPO}' must match owner/name shape (e.g. 'duncanjaffrey77-collab/staff-management')."
    continue
  fi

  if [[ ${first} -eq 0 ]]; then
    sleep "${STAGGER_SECONDS}"
  fi
  first=0

  start_runner "${NAME}" "${REPO}" "${LABELS}" "${TOKEN}" &
  CHILD_PIDS+=("$!")
done < "${TMP_RUNNERS}"

rm -f "${TMP_RUNNERS}"

# Graceful shutdown (v3.0.4). s6 sends SIGTERM to this script on container
# stop. Without a trap, bash exits and its background subshells running
# `./run.sh` die ungracefully, without flushing rotated credentials. The
# persisted `.credentials` then lags GitHub's server-side state, and on
# next restart the runner hits `PublicAccess … not authorized`.
#
# The trap forwards SIGTERM/SIGINT to all background subshells and waits
# up to GRACEFUL_TIMEOUT seconds for them to exit cleanly. actions/runner
# handles SIGTERM by flushing state and disconnecting cleanly from
# GitHub; 25s is comfortably above the few-seconds disconnect time and
# well below s6's own 30s kill-grace before SIGKILL.
GRACEFUL_TIMEOUT=25

shutdown_handler() {
  bashio::log.info "Shutdown signal received; forwarding to ${#CHILD_PIDS[@]} runner(s) (timeout ${GRACEFUL_TIMEOUT}s)."
  for pid in "${CHILD_PIDS[@]}"; do
    kill -TERM "${pid}" 2>/dev/null || true
  done
  # Wait up to GRACEFUL_TIMEOUT for all children to exit.
  local elapsed=0
  while [[ ${elapsed} -lt ${GRACEFUL_TIMEOUT} ]]; do
    local alive=0
    for pid in "${CHILD_PIDS[@]}"; do
      kill -0 "${pid}" 2>/dev/null && alive=$((alive + 1))
    done
    [[ ${alive} -eq 0 ]] && break
    sleep 1
    elapsed=$((elapsed + 1))
  done
  bashio::log.info "Graceful shutdown complete after ${elapsed}s."
  exit 0
}
trap shutdown_handler SIGTERM SIGINT

# Wait for ALL runners. If one dies the others keep working; we only exit
# (and let s6 restart) once everything is down.
wait
bashio::log.warning "All runners have exited. Service will restart."
